# sys-power/frequency-scaling/files/frequency-scaling.service
[Unit]
Description=Frequency Scaling

[Service]
ExecStart=/usr/bin/frequency-scaling
