# sys-power/frequency-scaling/files/frequency-scaling.service
[Unit]
Description=Check periodically that the proper scaling governor is used

[Service]
ExecStart=/usr/sbin/frequency-scaling
