# sys-power/laptop-tweaks/files/laptop-tweaks.service
[Unit]
Description= Tweak some power related laptop settings

[Service]
ExecStart=/usr/sbin/laptop-tweaks

[Install]
WantedBy=multi-user.target
